# crontab
* * * * *   su abuser -c "cd; flask update-feeds"
0 3 * * *   su abuser -c "cd; flask update-favicons"
