# crontab
* * * * *   su abuser -c "cd; flask update-feeds"
0 0 1 * *   su abuser -c "cd; flask update-favicons"
